# scripts/systemd/exodus-gw.target
[Unit]
Description=exodus-gw development services
